# donkeycron.crontab
* * * * * /usr/local/bin/jetsonvolt.bash >/tmp/jetsonvolt.log 2>&1
